# jetty.service
[Unit]
Description=Jetty Servlet Container

[Service]
ExecStart=/usr/share/jetty/jetty.sh

[Install]
WantedBy=multi-user.target
